# dist/packages/ovs-exporter/config/ovs-exporter.service
[Unit]
Description=OVS metrics exporter service

[Service]
ExecStart=<EXPORTER_PATH>/ovn-kube-util --loglevel=4 ovs-exporter --metrics-bind-address <LISTENING_IP>:9310

[Install]
WantedBy=multi-user.target
